# examples/recipes/07-system-service-providers/example-daemon/providers/init/runit/example-daemon/run
#!/bin/sh
exec 2>&1
exec chpst -u example-daemon /usr/bin/example-daemon \
  --config /etc/example-daemon/example-daemon.conf
